# .config/systemd/user/agentic-digest.service
[Unit]
Description=Daily agentic platform email digest via Claude agent

[Service]
Type=oneshot
ExecStart=/home/carter/scripts/run_agentic_digest.sh
TimeoutStartSec=600
