# .config/systemd/user/digests-daily.service
[Unit]
Description=Run all 5 daily email digests sequentially
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=/home/carter/scripts/run_all_digests.sh
TimeoutStartSec=21600
