# .config/systemd/user/gaming-digest.timer
[Unit]
Description=Fire gaming digest daily at 3pm ET

[Timer]
OnCalendar=*-*-* 19:00:00
Persistent=true

[Install]
WantedBy=timers.target
