# .config/systemd/user/homelab-steward.timer
[Unit]
Description=Daily homelab steward at 17:00 America/New_York (5pm ET)

[Timer]
OnCalendar=*-*-* 17:00:00 America/New_York
Persistent=true

[Install]
WantedBy=timers.target
