# .config/systemd/user/homelab-steward-resume.timer
[Unit]
Description=Post-boot resume trigger for homelab-steward

[Timer]
OnBootSec=2min

[Install]
WantedBy=timers.target
