# .config/systemd/user/homelab-backup-restore-drill.timer
[Unit]
Description=Monthly homelab backup restore drill

[Timer]
# 1st of each month at 12:00 UTC (07:00 ET). Runs after the 03:00 UTC backup
# has had time to settle on prior days; midday avoids contention with digests.
OnCalendar=*-*-01 12:00:00
Persistent=true

[Install]
WantedBy=timers.target