# .config/systemd/user/homelab-steward-notify.service
[Unit]
Description=Email failure notification for homelab-steward
After=homelab-steward.service
# Triggered by homelab-steward.service OnFailure=

[Service]
Type=oneshot
ExecStart=/home/carter/scripts/steward-notify-failure.sh
