# .config/systemd/user/hyperliquid-sdk.service
[Unit]
Description=Hyperliquid Ruby SDK autonomous maintenance run via omp (GLM 5.2)

[Service]
Type=oneshot
ExecStart=/home/carter/scripts/run_hyperliquid_sdk.sh
TimeoutStartSec=3600
