# .config/systemd/user/llama-server.service
[Unit]
Description=Local LLM Server (llama-swap proxy on gaming rig)
After=network-online.target
Wants=network-online.target

[Service]
Type=exec
# Unblock binary (Windows marks downloaded EXEs with Zone Identifier)
ExecStartPre=-/usr/bin/ssh -o ConnectTimeout=10 gamingrig "powershell -Command Unblock-File C:\\llm\\llama-swap.exe"
# Kill any stale processes before starting
ExecStartPre=-/usr/bin/ssh -o ConnectTimeout=10 gamingrig "taskkill /f /im llama-swap.exe"
ExecStartPre=-/usr/bin/ssh -o ConnectTimeout=10 gamingrig "taskkill /f /im llama-server.exe"
# Launch llama-swap via schtasks so it survives SSH disconnect
ExecStart=/usr/bin/ssh -o ServerAliveInterval=30 -o ServerAliveCountMax=3 gamingrig "C:\\llm\\llama-swap.exe --config C:\\llm\\config.yaml --listen 0.0.0.0:8080"
ExecStop=/usr/bin/ssh -o ConnectTimeout=10 gamingrig "taskkill /f /im llama-swap.exe"
Restart=on-failure
RestartSec=30
TimeoutStartSec=120
TimeoutStopSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
