# .config/systemd/user/opencode-homelab.service
[Unit]
Description=opencode web (homelab remote agent) - full-privilege, behind Cloudflare Access + OPENCODE_SERVER_PASSWORD
Documentation=https://opencode.ai/docs/
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/carter
Environment=HOME=/home/carter
# Full PATH so the agent can drive homelab ops (rbenv, fnm/node, kubectl, docker, git, etc.) - parity with the Claude daemon
Environment=PATH=/home/carter/.opencode/bin:/home/carter/.local/bin:/home/carter/.rbenv/bin:/home/carter/.rbenv/shims:/home/carter/.fnm:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
# OPENCODE_SERVER_PASSWORD (HTTP basic-auth second credential, enforced server-side) lives here, gitignored 600
EnvironmentFile=/home/carter/.config/opencode-homelab/env
# Bind loopback ONLY. The sole ingress is the CF tunnel (cloudflared runs on the host and reaches 127.0.0.1).
# NOT reachable on the LAN -> no bypass of Cloudflare Access. Mirrors the dependabot-webhook routing pattern.
ExecStart=/home/carter/.opencode/bin/opencode web --hostname 127.0.0.1 --port 48099
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
