# .config/systemd/user/pi-web-sessiond.service
[Unit]
Description=PI WEB session daemon

[Service]
Type=simple
ExecStart=/usr/bin/env /usr/bin/zsh -lc 'exec pi-web-sessiond'
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
