# .config/systemd/user/pi-web.service
[Unit]
Description=PI WEB server
After=pi-web-sessiond.service
Wants=pi-web-sessiond.service

[Service]
Type=simple
ExecStart=/usr/bin/env /usr/bin/zsh -lc 'exec pi-web-server'
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
