# .config/systemd/user/update-check.timer
[Unit]
Description=Fire homelab update audit nightly at 6am ET

[Timer]
OnCalendar=*-*-* 10:00:00
Persistent=true

[Install]
WantedBy=timers.target
