# .config/systemd/user/update-check.service
[Unit]
Description=Nightly homelab maintenance: auto-apply safe updates + validate + report

[Service]
Type=oneshot
ExecStart=/home/carter/scripts/run_update_check.sh
TimeoutStartSec=900
