# .config/systemd/user/update-check.service
[Unit]
Description=Nightly homelab maintenance: deterministic auto-apply + validate + report (1am ET)

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /home/carter/scripts/update_runner.py
TimeoutStartSec=900
